# deploy/systemd-user/auracore-backend.service
[Unit]
Description=AuraCore Backend (user)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__AURACORE_ROOT__/backend
EnvironmentFile=-__AURACORE_ROOT__/backend/.env
ExecStart=/bin/bash -lc 'exec "__AURACORE_ROOT__/scripts/run-backend-prod.sh"'
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
