# deploy/systemd-user/auracore-cloudflared.service
[Unit]
Description=AuraCore Cloudflare Tunnel (user)
After=network-online.target auracore-backend.service
Wants=network-online.target
Requires=auracore-backend.service

[Service]
Type=simple
WorkingDirectory=__AURACORE_ROOT__
ExecStart=/usr/bin/cloudflared --config %h/.cloudflared/auracore-config.yml tunnel run --token-file %h/.cloudflared/auracore-local-api.token
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
