# deploy/systemd-user/auracore-whatsapp-gateway.service
[Unit]
Description=AuraCore WhatsApp Gateway (user)
After=network-online.target auracore-backend.service
Wants=network-online.target
Requires=auracore-backend.service

[Service]
Type=simple
WorkingDirectory=/home/acer/Downloads/AuraCore/whatsapp-gateway
EnvironmentFile=-/home/acer/Downloads/AuraCore/whatsapp-gateway/.env
ExecStart=/bin/bash /home/acer/Downloads/AuraCore/scripts/run-gateway-prod.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
