# deploy/systemd-user/auracore-whatsapp-gateway.service
[Unit]
Description=AuraCore WhatsApp Gateway (user)
After=network-online.target auracore-backend.service
Wants=network-online.target auracore-backend.service

[Service]
Type=simple
WorkingDirectory=__AURACORE_ROOT__/whatsapp-gateway
EnvironmentFile=-__AURACORE_ROOT__/whatsapp-gateway/.env
ExecStart=/bin/bash -lc 'exec "__AURACORE_ROOT__/scripts/run-gateway-prod.sh"'
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
